# roles/netrom/files/systemd/netrom.target
[Unit]
Description=Netrom ports

[Install]
WantedBy=default.target
