# roles/direwolf/files/systemd/kissutil-cleanup.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=daily
