# units/ax25ports.target
[Unit]
Description=All AX25 ports are up

[Install]
WantedBy=radio.target
