# units/radio.target
[Unit]
Description=Radio services

[Install]
WantedBy=default.target
